# lnkact-monitor.service
[Unit]
Description=ASL3 Link Activity Monitor
After=asterisk.service
Requires=asterisk.service

[Service]
Type=simple
ExecStart=/etc/asterisk/scripts/lnkact-monitor/lnkact-monitor.sh
Restart=always
RestartSec=10
User=root
StandardOutput=journal
StandardError=journal
SyslogIdentifier=lnkact-monitor

[Install]
WantedBy=multi-user.target
